# systemd/clawhip-project-executor-keepalive@.service
[Unit]
Description=Project executor keepalive for %i
After=clawhip.service
Requires=clawhip.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/mei/.config/ignition/bin/projectctl.py keepalive %i executor --timeout 5
